# build/device/zero2w/device/rootfs-overlay/etc/systemd/system/pirate-wifi.service
[Unit]
Description=Apply Pirate WiFi config from /config/wifi.cfg
ConditionPathExists=/boot/firstboot_done
Wants=network.target
After=NetworkManager.service network.target

[Service]
Type=oneshot
ExecStart=/usr/local/pirate/wifi.sh /config/wifi.cfg
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
