# build/device/zero2w/device/rootfs-overlay/etc/systemd/system/pirate-firstboot.service
[Unit]
Description=Perform Pirate first-boot commands
ConditionPathExists=!/boot/firstboot_done
After=local-fs.target
DefaultDependencies=yes

[Service]
Type=oneshot
ExecStart=/usr/local/pirate/first-boot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
